# sanmar_integration/deploy/systemd/sanmar-api.service
[Unit]
Description=SanMar local cache API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sanmar
WorkingDirectory=/opt/sanmar
EnvironmentFile=/opt/sanmar/.env
ExecStart=/opt/sanmar/.venv/bin/python -m sanmar serve-api --port 8000
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
